# files/target/linux/ath79/mikrotik/base-files/etc/uci-defaults/99-bootbeep
#!/bin/sh
# Enable the startup beep service on first boot (persists in overlay,
# so it then runs on every subsequent boot).
[ -x /etc/init.d/bootbeep ] && /etc/init.d/bootbeep enable
exit 0
